# etc/udev/rules.d/66-input-optimization.rules
# ============================================================================
# Input Device Optimization
# ============================================================================
# Only rules with verifiable sysfs effect are included.
# Keyboard repeat rate, mouse acceleration, and libinput options must be
# configured in X11 (xorg.conf.d), Wayland compositor config, or user
# session scripts — udev has no mechanism to inject these settings.

# USB HID devices (keyboards, mice, game controllers, drawing tablets):
# disable autosuspend on the parent USB device node.
# power/autosuspend_delay_ms lives on the device node (1-7), not the
# interface node (1-7:1.0), so we match the interface and write to its parent.
ACTION=="add", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_interface", \
  ATTRS{bInterfaceClass}=="03", \
  RUN+="/bin/sh -c 'echo -1 > /sys$(dirname %p)/power/autosuspend_delay_ms'", \
  RUN+="/bin/sh -c 'echo on > /sys$(dirname %p)/power/control'"

# Removed rules and reasons:
# - echo > /proc/sys/kernel/autorepeat_interval: path does not exist
# - ENV{LIBINPUT_IGNORE_DUPLICATE_EVENT}: not a real libinput udev property
# - ENV{MOUSE_ACCEL}/ENV{MOUSE_SPEED}: not recognized by X11 or libinput
# - echo 1000 > /sys/module/usbhid/parameters/jitter: 1000ms jitter buffer
#   increases lag rather than reducing it, and the path may not exist
