# etc/systemd/system/systemd-logind.service.d/10-kiro-session.conf
# ============================================================================
# Kiro-ISO Logind Configuration (Session Management)
# ============================================================================
# Optimizes session behavior for responsiveness

[Service]

# Timeout for logind operations
# Recommendation: 10s (responsive login)
TimeoutStopSec=10s

# Restart on failure
# Default: no
# Recommendation: yes (ensures session management works)
Restart=on-failure
RestartSec=5s

# Kill all session processes on logout
# Recommendation: yes (clean exit)
KillMode=mixed
KillSignal=SIGTERM
